# templates/dnsmasq@.service.j2
# dnsmasq systemd service template
# Generated by Ansible - DO NOT EDIT MANUALLY

[Unit]
Description=dnsmasq - A lightweight DHCP and caching DNS server (%i)
After=network.target
Wants=nss-lookup.target
Before=nss-lookup.target

[Service]
Type=forking
PIDFile=/run/dnsmasq-%i.pid
ExecStartPre=/usr/sbin/dnsmasq --test --conf-dir=/etc/dnsmasq.d/%i,*.conf
ExecStart=/usr/sbin/dnsmasq --conf-dir=/etc/dnsmasq.d/%i,*.conf --pid-file=/run/dnsmasq-%i.pid
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
